# sbctl.service
[Unit]
Description=sbctl - sing-box 进程级代理控制器
Documentation=https://github.com/decker502/sbctl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sbctl
EnvironmentFile=/opt/sbctl/.env
ExecStart=/bin/bash -c 'source /opt/sbctl/.env && exec /opt/sbctl/sbctl tun "$PROXY_NAME" $PROCESS_NAMES'
ExecStopPost=/opt/sbctl/clean.sh
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

# 安全设置
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
